# config/tlgbot.conf
[program:tlgbot]
command=/home/www/kakula_bot/start_tlgbot.sh
user=www
process_name=%(program_name)s
numproc=1
autostart=1
autorestart=1
redirect_stderr=true